# puppet/modules/systembase_cdn_qa/templates/monit.conf
# system monit 

start on runlevel [35]
limit nofile 1048576 1048576
limit nproc  8192 8192
stop on runlevel [!$RUNLEVEL]

respawn
exec /usr/bin/monit -I
